# adoctor-diag-executor/adoctor-diag-executor.service
[Unit]
Description=A-doctor's Diagnose Executor
After=network.target

[Service]
Restart=on-failure
RestartSec=5
ExecStart=/usr/bin/adoctor-diag-executor start
ExecStop=/usr/bin/adoctor-diag-executor stop

[Install]
WantedBy=multi-user.target
